# etc/systemd/system/connman.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/connmand -n -r
