# etc/systemd/system/connman.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/connmand --nodaemon --nodnsproxy
