# home/waldrep/.config/systemd/user/keybase.gui.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/electron13 --enable-features=UseOzonePlatform --ozone-platform=wayland /usr/share/keybase-app
