# etc/systemd/system/wpa_supplicant.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/wpa_supplicant -u -O /var/run/wpa_supplicant
